# infra/server/systemd/virtualwait-gateway.service
[Unit]
Description=VirtualWait signed verification gateway
After=network.target

[Service]
Type=simple
User=virtualwait
Group=virtualwait
WorkingDirectory=/opt/virtualwait/services/sdgb-gateway
EnvironmentFile=/etc/virtualwait/gateway.env
Environment=PYTHONPATH=/opt/virtualwait/services/sdgb-gateway/src
ExecStart=/usr/bin/python3 -m virtualwait_gateway
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/virtualwait/gateway

[Install]
WantedBy=multi-user.target
