# infra/server/systemd/virtualwait-maintenance.service
[Unit]
Description=VirtualWait queue and transient-data maintenance worker
After=network.target virtualwait-web.service

[Service]
Type=simple
User=virtualwait
Group=virtualwait
WorkingDirectory=/opt/virtualwait/apps/web
EnvironmentFile=/etc/virtualwait/web.env
ExecStart=/usr/bin/npm run maintenance
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/virtualwait/web
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
